# resources/sidebyside/resmon2-manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='resmon2'>

<service
	name='site/resmon2'
	type='service'
	version='1'>

	<create_default_instance enabled='false' />

	<single_instance />

	<dependency
		name='config-file'
		grouping='require_all'
		restart_on='refresh'
		type='path'>
			<service_fmri
				value='file://localhost/opt/resmon2/resmon2.conf' />
	</dependency>

	<dependency
		name='resmon2_multi-user-server'
		grouping='require_all'
		restart_on='none'
		type='service'>
		<service_fmri value='svc:/milestone/multi-user-server' />
	</dependency>

	<exec_method
	    type='method'
	    name='start'
	    exec='/opt/resmon2/resmon -c /opt/resmon2/resmon2.conf'
	    timeout_seconds='60' />

	<exec_method
	    type='method'
	    name='stop'
	    exec=':kill'
	    timeout_seconds='60' />

	<exec_method
	    type='method'
	    name='refresh'
	    exec=':kill -HUP'
	    timeout_seconds='60' />

	<stability value='Unstable' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			Resmon
			</loctext>
		</common_name>
	</template>
</service>

</service_bundle>
